# templates/supervisor.ini
[program:{{ app_name }}]
command=/usr/local/bin/uwsgi
  --socket /tmp/{{ app_name }}.sock
  --logto {{ webapps_dir }}/{{ app_name }}/log/application.log
  --home {{ webapps_dir }}/{{ app_name }}/venv
  --pythonpath {{ webapps_dir }}/{{ app_name }}/src
  --wsgi-file {{ webapps_dir }}/{{ app_name }}/src/{{ wsgi_file }}
  --callable {{ wsgi_callable }}
  --max-requests 1000
  --master
  --processes 1
  --chmod
directory={{ webapps_dir }}/{{ app_name }}/src
autostart=true
autorestart=true